# zram-config.service
[Unit]
Description=zram-config
DefaultDependencies=no
After=local-fs.target
Before=systemd-journal-flush.service

[Service]
Type=oneshot
ExecStart= /usr/local/bin/zram-config start
ExecStop= /usr/local/bin/zram-config stop
TimeoutStartSec=120
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target

