# zram-config.service
[Unit]
Description=zram-config
Documentation=https://github.com/ecdye/zram-config/blob/main/README.md
Requires=local-fs.target
After=local-fs.target

[Service]
Type=exec
TimeoutSec=120
RemainAfterExit=yes
ExecStart=/usr/local/sbin/zram-config "start"
ExecReload=/usr/local/sbin/zram-config "start"
ExecStop=/usr/local/sbin/zram-config "stop"

[Install]
WantedBy=multi-user.target
