# zsync.cron
#!/usr/bin/env bash

/usr/sbin/zram-config "sync"
