# tests/rc.local
cd /opt/zram || exit 1
./install.bash
zramctl --output-all
shutdown -h now
